# trash/.config/systemd/user/empty-trash.service
[Unit]
Description=Empty trash

[Service]
Type=oneshot
ExecStart=/usr/bin/trash-empty -f 30
